# remote/sys/crontab
0 * * * * /root/IPT400S/remote/db/sync.sh
5 0 * * * /root/IPT400S/remote/copyclean.sh

